# files/ansible/crontab
SHELL=/bin/sh
HOME=/srv/ansible
*/5 * * * * ansible ansible-playbook playbooks/ansible/schedule/update_duckdns.yml --vault-pass-file keys/auto_vault_pw >> /var/log/ansible/duckdns.log